# root/etc/services.d/dnscrypt-proxy/run
#!/usr/bin/with-contenv sh

echo "Starting dnscrypt-proxy..."
/dnscrypt/dnscrypt-proxy -config /etc/dnscrypt-proxy.toml
