# root/etc/services.d/socks/run
#!/usr/bin/with-contenv sh

echo "Starting socks proxy..."
s6-setuidgid abc socks -t
